# scripts/adsb-aamt-feed.service

[Unit]
Description=adsb-aamt-feed
Wants=network.target
After=network.target

[Service]
User=adsb-aamt
RuntimeDirectory=adsb-aamt-feed
ExecStart=/usr/local/share/adsb-aamt/adsb-aamt-feed.sh
Type=simple
Restart=always
RestartSec=30
StartLimitInterval=1
StartLimitBurst=100
SyslogIdentifier=adsb-aamt-feed
Nice=-1

[Install]
WantedBy=default.target
